# systemd/sserveros-agent-api.service.in
[Unit]
Description=sserveros node Agent API
PartOf=sserveros-controller.target sserveros-agent.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__SCRIPT_DIR__
EnvironmentFile=-__SCRIPT_DIR__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON_BIN__ __SCRIPT_DIR__/agent_api.py
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=15

[Install]
WantedBy=sserveros-controller.target sserveros-agent.target
